# _tmp_lokalny_worker.service
[Unit]
Description=Lokalny Obiektyw background job worker
After=network.target lokalny_obiektyw.service
Requires=lokalny_obiektyw.service

[Service]
Type=simple
User=debian
WorkingDirectory=/home/debian/lokalny_obiektyw
EnvironmentFile=-/etc/lokalny_obiektyw/lokalny_obiektyw.env
Environment=PYTHONPATH=/home/debian/lokalny_obiektyw
ExecStart=/home/debian/lokalny_obiektyw/venv/bin/python /home/debian/lokalny_obiektyw/backend/worker.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target